# screen-time.service
[Unit]
Description=Monitor computer usage

[Service]
Type=oneshot
ExecStart=/usr/bin/screen-time_startup.sh

[Install]
WantedBy=multi-user.target