# screen-time-daemon.service
[Unit]
Description=Monitor screen on time.

[Service]
Type=oneshot
ExecStart=/usr/bin/screen-time-daemon

[Install]
WantedBy=multi-user.target