# screen-time.service
[Unit]
Description=Monitor screen on time.

[Service]
Type=oneshot
ExecStart=/usr/bin/screen-time

[Install]
WantedBy=multi-user.target